# telegram-bot.service.example
[Unit]
Description=Telegram Chat Analyzer Bot
Documentation=https://github.com/your-username/ChatSum
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=your_username
Group=your_username
WorkingDirectory=/home/your_username/ChatSum
ExecStart=/usr/bin/python3 /home/your_username/ChatSum/main.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Environment
Environment="PYTHONUNBUFFERED=1"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bot

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/your_username/ChatSum

[Install]
WantedBy=multi-user.target

